# pipeline/crontab
# Chicago NDVI Pipeline — Monthly cron schedule
# Runs on the 1st of each month at 06:00 UTC
# Processes the prior month's satellite imagery
SHELL=/bin/bash
PATH=/usr/local/sbin:/usr/local/bin:/sbin:/bin:/usr/sbin:/usr/bin

0 6 1 * * root cd /app && python pipeline.py --month $(date -d "last month" +\%Y-\%m) >> /var/log/pipeline.log 2>&1
